# crontab
* * * * * php /var/apps/vdt/xci/jobby.php >> /var/apps/vdt/xci/logs/jobby.log 2>&1